# contrib/crontab
# Run strava-backup every Monday at 00:00
0 0 * * 1 /usr/local/bin/strava-backup --quiet
