# contrib/strava-backup@.service
[Unit]
Description=Run strava-backup for user %i

[Service]
User=%i
ExecStart=/usr/local/bin/strava-backup --quiet
